# flo-pes-scheduler/config/supervisord.conf
; Configuration docs:
; http://supervisord.org/configuration.html
; https://gist.github.com/tsabat/1528270

[unix_http_server]
file = /tmp/supervisor.sock

[inet_http_server]
port = 9001

[supervisord]
logfile = /var/logs/supervisor/supervisord.log
childlogdir = /var/logs/supervisor
logfile_maxbytes = 10MB
logfile_backups = 10
loglevel = info
nocleanup = true
nodaemon = true
pidfile = /tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:flo-pes-scheduler]
environment=PATH="/opt/conda/envs/fpw/bin"
command=celery beat -A sensor_api.celery -l INFO
directory = /app/celery
numprocs=1
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
autostart = true
autorestart = true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs=600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

[program:worker-ping]
environment=PATH="/opt/conda/envs/fpw/bin"
command = python app.py
directory = /app
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
